# bjx_power/deploy/bjx-daily.timer
# 用户级 systemd timer —— 与 bjx-daily.service 一起放入 ~/.config/systemd/user/
[Unit]
Description=BJX daily crawl timer (weekday mornings, Asia/Shanghai)

[Timer]
# 服务器时区已设为 Asia/Shanghai(timedatectl), 此为北京时间周一至周五 08:30
OnCalendar=Mon..Fri 08:30:00
Persistent=true

[Install]
WantedBy=timers.target
